# systemd/aero-halo-360-health.timer
[Unit]
Description=жЇЏ 5 з§’жЈЂжџҐ AeroHalo360 дј ж„џдёЋеЏ‘йЂЃеЃҐеє·зЉ¶жЂЃ

[Timer]
OnBootSec=20
OnUnitActiveSec=5
AccuracySec=1
Unit=aero-halo-360-health.service

[Install]
WantedBy=timers.target
